# configs/systemd/recepio-backup.timer
[Unit]
Description=Recepio Secrets Backup Timer

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
